# supervisord.conf
[program:metricbot]
command=node ./src/bot.js    ; the program (relative uses PATH, can take args)
autostart=true               ; start at supervisord start (default: true)
autorestart=true             ; when to restart if exited after running (def: unexpected)
directory=/var/www/html/metric_units_bot
